# systemd/full/user/pipewire-media-session.service
[Service]
NoNewPrivileges=no
MemoryDenyWriteExecute=no
LockPersonality=no
RestrictNamespaces=no
