# systemd/haveged.service
[Unit]
After=apparmor.service systemd-tmpfiles-setup-dev.service
