# systemd/systemd-udevd.service
[Unit]
After=apparmor.service systemd-sysusers.service systemd-hwdb-update.service