# systemd/systemd-networkd.service
[Unit]
After=apparmor.service systemd-udevd.service network-pre.target systemd-sysusers.service systemd-sysctl.service
