# systemd/systemd-journald.service
[Unit]
After=apparmor.service systemd-journald.socket systemd-journald-dev-log.socket systemd-journald-audit.socket syslog.socket
