# systemd/systemd-timesyncd.service
[Unit]
After=apparmor.service systemd-sysusers.service
