# systemd/full/system/archlinux-keyring-wkd-sync.service
[Service]
AppArmorProfile=&archlinux-keyring-wkd-sync
