# systemd/full/system/cloud-init-hotplugd.service
[Service]
AppArmorProfile=&cloud-init-hotplugd.service
