# systemd/full/system/debug-shell.service
[Service]
AppArmorProfile=debug-shell.service