# systemd/full/system/fwupd-refresh.service
[Service]
ProtectKernelModules=no
RestrictRealtime=no
ProtectKernelModules=no
