# systemd/full/system/geoclue.service
[Service]
NoNewPrivileges=no
MemoryDenyWriteExecute=no
ProtectKernelTunables=no
ProtectKernelModules=no
RestrictRealtime=no
