# systemd/full/system/grub-common.service
[Service]
AppArmorProfile=grub-common.service