# systemd/full/system/man-db.service
[Service]
AppArmorProfile=man-db.service