# systemd/full/system/nm-priv-helper.service
[Service]
AppArmorProfile=&nm-priv-helper
