# systemd/full/system/power-profiles-daemon.service
[Service]
AppArmorProfile=&power-profiles-daemon