# systemd/full/system/secureboot-db.service
[Service]
AppArmorProfile=secureboot-db.service
