# systemd/full/system/snapd.system-shutdown.service
[Service]
AppArmorProfile=snapd.system-shutdown.service