# systemd/full/system/system-update-cleanup.service
[Service]
AppArmorProfile=system-update-cleanup.service