# systemd/full/system/systemd-journal-remote.service
[Service]
AppArmorProfile=&systemd-journal-remote