# systemd/full/system/systemd-oomd.service
[Service]
NoNewPrivileges=no
ProtectClock=no
