# systemd/full/system/usb_modeswitch@.service
[Service]
AppArmorProfile=usb_modeswitch.service