# systemd/default/system/dbus-broker.service
[Unit]
After=apparmor.service

[Service]
AppArmorProfile=dbus-system
